# void/apps/onedrive/service/log/run
#!/bin/sh
exec svlogd -t /etc/sv/onedrive/log
