# void/apps/nordvpn/service/log/run
#!/bin/sh
exec svlogd -t /etc/sv/nordvpn/log
